# get-public-ip/get-public-ip.cron
#!/usr/bin/env bash
#
# Periodically check public IP
# sends warnings if change is detected

get-public-ip.sh >/dev/null 2>&1
